# main/src/main/dist/init/freki-web.conf

description     "The Freki sensor data query HTTP service"

start on runlevel [345]
stop on runlevel [0126]

chdir /opt/freki

limit nofile 8192 8192

respawn
respawn limit 5 60

kill timeout 20

exec /opt/freki/bin/freki web
